# provision/fluxbox.conf
description "fluxbox"

start on started xvfb
stop on stopping xvfb

respawn
respawn limit 5 30

env DISPLAY=:99.0
exec su vagrant -c "/usr/bin/fluxbox"
