# provision/skype.conf
description "skype"

start on started fluxbox
stop on stopping fluxbox

respawn
respawn limit 5 30

env DISPLAY=:99.0
exec su vagrant -c "/bin/bash /vagrant/provision/run-skype.sh"
