# provision/x11vnc.conf
description "x11vnc"

start on started xvfb
stop on stopping xvfb

respawn
respawn limit 5 30

env PASSWORD=changeme
exec /usr/bin/x11vnc -display :99 -passwd "$PASSWORD" -forever
